# fermatai-slow-claude.timer
[Unit]
Description=Slow Claude monitor timer (saatte 1)

[Timer]
OnCalendar=*:00:00
Persistent=true
Unit=fermatai-slow-claude.service

[Install]
WantedBy=timers.target
