# eyotek_agent/systemd/fermat-chrome-cdp.service
[Unit]
Description=Fermat Chrome CDP — Eyotek 7/24 baglanti icin headless Chromium
Documentation=https://github.com/zekigoksal-source/fermatai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=neo
Group=neo
Environment="HOME=/home/neo"
Environment="DISPLAY="
# Chrome CDP — port 9333, headless new, persistent profile.
# ExecStart bir wrapper script'e delege edildi (26 May 2026): Playwright guncellemesi
# chromium binary klasorunu degistiriyor (chromium-NNNN), eskiden hardcode oldugu icin
# 203/EXEC flapping yapiyordu. Wrapper en guncel chromium-* build'i otomatik bulur.
ExecStartPre=/bin/mkdir -p /home/neo/.fermat-chrome /home/neo/.fermat-chrome-logs
ExecStart=/opt/fermatai/eyotek_agent/chrome_cdp_start.sh

# Restart on failure (network glitch, OOM, etc)
Restart=always
RestartSec=5
# Memory limit — Chrome'u 1GB ile kapatma (idle ~200MB)
MemoryMax=1G
# CPU quota (4 worker bridge ile race olmasin)
CPUQuota=50%

# Stdout/stderr -> journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fermat-chrome-cdp

[Install]
WantedBy=multi-user.target
