# vps_setup/systemd/fermatai-smart-sync.service
[Unit]
Description=FermatAI - Smart sync (Eyotek delta sınav verisi)
Documentation=Oturum 25.9 T6
After=fermatai-eyotek-daily.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=neo
Group=neo
WorkingDirectory=/opt/fermatai/eyotek_agent
Environment="PATH=/opt/fermatai/.venv/bin:/usr/bin"
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/fermatai/.playwright-browsers
EnvironmentFile=/opt/fermatai/.env

# smart_sync.py: incremental — sınav sayısı değişmemişse skip eder
# resume modu: kaldığı yerden devam
ExecStart=/opt/fermatai/.venv/bin/python smart_sync.py --resume

StandardOutput=journal
StandardError=journal
SyslogIdentifier=fermatai-smart-sync
TimeoutStartSec=3600
