# vps_setup/systemd/fermatai-backup.timer
[Unit]
Description=FermatAI - Daily backup timer (03:00 UTC)
Documentation=Oturum 25.9 T5

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true
Unit=fermatai-backup.service

[Install]
WantedBy=timers.target
